# deploy/rpm/infinit.e-record-engine.rpm/SOURCES/infinit.e-record-engine/etc/cron.d/infinite-logstash
# /etc/cron.d/infinite: infinite crontab
# Unlike any other crontab you don't have to run the `crontab'
# command to install the new version when you edit this file
# since it is in /etc/cron.d. These files also have username fields,
# that none of the other crontabs do.

# use full cron commands after a cron modification
# stop the cron service
# sudo stop cron
# start the cron service
# sudo start cron
# verify the cron log
# tail -f /var/log/syslog
# verify the infinite log 
# tail -f /opt/infinite-home/logs/logfile.log

SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
MAILTO=""

# m h dom mon dow user  command

* * * * * root	if [ -f /opt/logstash-infinite/RESTART_LOGSTASH ]; then service logstash restart && rm -f /opt/logstash-infinite/RESTART_LOGSTASH; fi

*/10 * * * * root	sh /opt/logstash-infinite/scripts/remove_logstash_tmpfiles.sh

* * * * * root /sbin/service logstash status >> /dev/null; RETCODE=$?; if [ $(find /etc/logstash/conf.d/ -name "*.conf"|wc -l) -gt 0 ] && [ $RETCODE -eq 2 ] ; then touch /opt/logstash-infinite/RESTART_LOGSTASH; fi

#(EOL defender!)
