# deploy/rpm/infinit.e-db-instance.rpm/SOURCES/infinit.e-db-instance/etc/cron.d/infinite-db
# /etc/cron.d/infinite: infinite crontab
# Unlike any other crontab you don't have to run the `crontab'
# command to install the new version when you edit this file
# since it is in /etc/cron.d. These files also have username fields,
# that none of the other crontabs do.

# use full cron commands after a cron modification
# stop the cron service
# sudo stop cron
# start the cron service
# sudo start cron
# verify the cron log
# tail -f /var/log/syslog
# verify the infinite log 
# tail -f /opt/infinite-home/logs/logfile.log

SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

# m h dom mon dow user  command

# Backup (works out what to do based on config)
# (Fri night backup handled by the compaction scripts below)
# (Don't backup on Sat night because of the big batch script that gets run)
01 01    * * 0-5   root    /opt/db-home/master_backup-script.sh > /tmp/backup_script.log

# Fri night: compaction time!
# First the secondaries, then the master a few hours laster
05 00	* * 6	root	/opt/db-home/master_compact-script.sh 1	> /tmp/compact_script.log
05 03	* * 6	root	/opt/db-home/master_compact-script.sh 2 >> /tmp/compact_script.log

# For the (MT) midnight one, do more - used to rotate the mongodb logs
00 01    * * *   root    /opt/db-home/rotate-logs-script.sh

#Every 5 minutes: check if we need to replace a temporary instance-arbiter with a permanent config arbiter
*/5 * * * * root sh /opt/db-home/arbiter-manager.sh > /dev/null

#(EOL defender!)
