# deploy/config_management_rpm/interface-engine.rpm/SOURCES/ikanow-interface-engine/etc/cron.d/tomcat-interface-engine
# /etc/cron.d/infinite: infinite crontab
# Unlike any other crontab you don't have to run the `crontab'
# command to install the new version when you edit this file
# since it is in /etc/cron.d. These files also have username fields,
# that none of the other crontabs do.

# use full cron commands after a cron modification
# stop the cron service
# sudo stop cron
# start the cron service
# sudo start cron
# verify the cron log
# tail -f /var/log/syslog
# verify the infinite log 
# tail -f /opt/infinite-home/logs/logfile.log

SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
MAILTO=""

#Every minute, check the interface engine is still running, restart it if not
* * * * * root	if ! pgrep -u tomcat -f org.apache.catalina.startup.Bootstrap > /dev/null && [ -f /var/lock/subsys/tomcat-interface-engine ]; then service tomcat-interface-engine start; fi 

#(EOL defender!)
