# deploy/rpm/infinit.e-processing-engine.rpm/SOURCES/infinit.e-processing-engine/etc/cron.d/infinite-px-engine
# /etc/cron.d/infinite: infinite crontab
# Unlike any other crontab you don't have to run the `crontab'
# command to install the new version when you edit this file
# since it is in /etc/cron.d. These files also have username fields,
# that none of the other crontabs do.

# use full cron commands after a cron modification
# stop the cron service
# sudo stop cron
# start the cron service
# sudo start cron
# verify the cron log
# tail -f /var/log/syslog
# verify the infinite log 
# tail -f /opt/infinite-home/logs/logfile.log

SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

# m h dom mon dow user  command

# used to check all instances still alive
* * * * * root	if [ ! -f /opt/infinite-home/bin/STOPFILE ]; then service infinite-px-engine watchdog; fi
# Also, regulate index size in the event of node failure (fails silently if API not up, which is fine):
* * * * * root	curl -s 'http://localhost:8080/api/auth/login/ping/ping' > /dev/null

#Sync DB+index every hour
00 *  * * * root touch /opt/infinite-home/bin/SYNC_FILE

#Reset sources every 24 hours (use script version rather than programmatic version)
#00 00 * * * root touch /opt/infinite-home/bin/RESET_FILE
00 00 * * * root /opt/infinite-home/bin/reset_bad_harvest.sh

# Saturday night: big resync! (takes about 30 minutes, only 1 harvester runs it but all stop - mainly managed by db.feature.sync_lock)
30 01	 * * 7	 root	if [ ! -f /opt/infinite-home/bin/STOP_BATCH_SYNC_FILE ]; then service infinite-px-engine stop; fi
00 02    * * 7   root   if [ ! -f /opt/infinite-home/bin/STOP_BATCH_SYNC_FILE ]; then /opt/infinite-home/bin/sync_features.sh; fi
30 02	 * * 7	 root	if [ ! -f /opt/infinite-home/bin/STOP_BATCH_SYNC_FILE ]; then service infinite-px-engine start; fi

# Sunday afternoon: weekly sources report
00 14	 * * 0	 root	 /opt/infinite-home/bin/weekly_sources_report.sh

#Custom Mapreduce check/run jobs every 1 minutes
*/1 * * * * tomcat	if [ ! -f /opt/infinite-home/bin/STOP_CUSTOM ]; then /opt/infinite-home/bin/custommr.sh; fi
